# edge-gateway/services/rs232-gateway/systemd/rs232-gateway.service
[Unit]
Description=EdgeVib RS232 Gateway — Serial-to-MQTT Bridge (Path B Backup)
Documentation=https://github.com/EnterWorldDoor/smart-vibration-monitoring-system
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=orangepi
Group=dialout
WorkingDirectory=/opt/edge-gateway
ExecStart=/usr/local/bin/rs232-gateway /opt/edge-gateway/config/rs232-gateway.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rs232-gateway

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ReadWritePaths=/dev /var/log
ReadOnlyPaths=/opt/edge-gateway/config
MemoryMax=50M

# Resource limits
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
